# crontab
# Crontab file for supercronic

# Run every hour

@hourly python manage.py runcrons
